# etc/cron.d/teamcity
40 8 * * * root /etc/init.d/teamcity restart > /dev/null 2>&1
