# debian/teamcity-server.postinst
#!/bin/sh -e
# postinst script for teamcity-server
set -e

U="teamcity"

case $1 in
    configure)

        #DEBHELPER#

        chown $U /var/lib/teamcity
        chown $U /var/lib/teamcity/BuildServer
        chown $U /var/lib/teamcity/BuildServer/config
        chown $U /var/lib/teamcity/BuildServer/config/*.properties
        chown $U /var/lib/teamcity/BuildServer/lib
        chown $U /var/lib/teamcity/BuildServer/lib/jdbc
        chown $U /var/lib/teamcity/BuildServer/lib/jdbc/*.jar

        mysql -e "GRANT ALL PRIVILEGES ON teamcity.* TO 'teamcity'@'localhost' IDENTIFIED BY 'teamcity'" 2>/dev/null || true
        mysql -uteamcity -pteamcity -e "CREATE DATABASE teamcity DEFAULT CHARSET utf8" 2>/dev/null || true

        /etc/init.d/teamcity upgrade @TEAMCITY_VERSION@

        update-rc.d teamcity defaults 2>/dev/null || true

        ln -sf /etc/nginx/sites-available/teamcity /etc/nginx/sites-enabled/teamcity 2>/dev/null || true
        ln -sf /etc/nginx/sites-available/teamcity-ssl /etc/nginx/sites-enabled/teamcity-ssl 2>/dev/null || true
        if [ -f /etc/nginx/nginx-default.conf ]; then
            nginx -tc /etc/nginx/nginx-default.conf 2>/dev/null || true
        fi
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        ;;

esac

exit 0
